# local/etc/profile.d/20-path.sh
# Standard locations for searching for executables

# Add user path:
path prepend "$(readlink -f "$HOME/bin")"

