# local/etc/profile.d/60-editor.sh
# Text editors/pagers

export EDITOR='/usr/bin/nvim'
export VISUAL='/usr/bin/nvim'
export PAGER='/usr/bin/less'
